# dist/sma.cron
* * * * * /opt/sbfspot/SBFspot -v -cfg/config/SBFspot.cfg > /dev/null
